# image/pd.conf
[program:pd]
command=/bin/pd-server --client-urls http://127.0.0.1:2379 --data-dir "/data/pd"
