# image/tidb.conf
[program:tidb]
command=/bin/tidb-server --config "/conf/tidb.toml" -P 4000 --store tikv --path 127.0.0.1:2379 --host 0.0.0.0 --advertise-address 127.0.0.1