# image/tikv.conf
[program:tikv]
command=/bin/tikv-server --status-addr 127.0.0.1:20180  --pd 127.0.0.1:2379 -A 127.0.0.1:20160 -C "/conf/tikv.toml" -s "/data/tikv"
